# plasma-i3.service
[Install]
WantedBy=plasma-workspace.target

[Unit]
Description=i3wm
Before=plasma-workspace.target

[Service]
ExecStart=/usr/bin/i3
Slice=session.slice
Restart=on-failure
